# infra/workflow.cron
# Run workflow.py every minute for testing
0 3 * * * cd /app && /usr/local/bin/python scripts/workflow.py >> /var/log/cron.log 2>&1
